# supervisor_config_files/vnmk-server.conf.example
[program:vnmkserver]
command=python3 -m vnmk.server <CONFIG YAML PATH>.yaml
user=nobody
autostart=true
autorestart=true
